# conf/supervisord/coinboot.ini
[supervisord]
nodaemon=true

[program:dnsmasq]
command=/usr/sbin/dnsmasq --no-daemon
stdout_logfile=/var/log//%(program_name)s.log
stderr_logfile=/var/log//%(program_name)s.err
username=root
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=900
stdout_logfile= /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
username=www-data
autorestart=true
